# contrib/supervisord/descartes-bi.conf
[program:descartes-bi]
command = /usr/share/descartes-bi/contrib/scripts/gunicorn_start.sh
user = www-data
redirect_stderr = true
